# crontabs/root
*/5 * * * * /usr/local/bin/docker_volume_cleanup.py >>/var/log/docker_volume_cleanup.log 2>&1
